# etc/supervisor.d/caddy.conf
[program:caddy]
command=/usr/bin/caddy -conf="/etc/Caddyfile" %(ENV_CADDY_OPTIONS)s
stdout_logfile=/var/log/caddy.log
stderr_logfile=/var/log/caddyerr.log
