# etc/supervisor.d/caddy.conf
[program:caddy]
command=caddy --conf /etc/Caddyfile %(ENV_CADDY_OPTIONS)s
