# etc/supervisor.d/docker-gen.conf
[program:docker-gen]
command=docker-gen -notify "supervisorctl restart caddy" -watch /etc/caddy-proxy/caddy.tmpl /etc/Caddyfile
