# etc/supervisor.d/docker-gen.conf
[program:docker-gen]
command=docker-gen -notify "supervisorctl signal USR1 caddy" -watch /etc/caddy-proxy/caddy.tmpl /etc/Caddyfile
